# list/location.list
/home/ari/.bashrc -> dotfiles/shells/bash
/home/ari/.profile -> dotfiles/shells/bash
/home/ari/.config -> dotfiles/config
/home/ari/Documents/dwm_bar.sh -> dotfiles/suckless
/home/ari/Suckless -> dotfiles/suckless
/home/ari/.gtkrc-2.0 -> dotfiles/etc
/home/ari/Media/Pictures/wallpaper.png -> dotfiles/etc
/home/ari/.mailcap -> dotfiles/etc
/home/ari/.mutt/muttrc -> dotfiles/etc
/home/ari/.lynxrc -> dotfiles/etc
/home/ari/.gitconfig -> dotfiles/etc
/home/ari/.mozilla/firefox/3s4h1qq0.default-release -> dotfiles/etc/firefox
/home/ari/.icons -> dotfiles/icons
/home/ari/.scripts -> dotfiles/custom
/home/ari/.vim -> dotfiles/editors/vim
/home/ari/.vimrc -> dotfiles/editors/vim
/home/ari/.config/VSCodium/User/snippets -> dotfiles/programming/VScodium
/home/ari/.xinitrc -> dotfiles/core
/home/ari/.xprofile -> dotfiles/core
/etc/default/grub -> dotfiles/core
/etc/hosts -> dotfiles/core
/etc/dracut.conf -> dotfiles/core
/etc/dracut.conf.d -> dotfiles/core
/etc/sudoers -> dotfiles/core
/etc/X11/xorg.conf.d -> dotfiles/core
/etc/rc.conf -> dotfiles/core
/etc/modprobe.d/kernel_unfreeze_rtw.conf -> dotfiles/core
/etc/conf.d -> dotfiles/core
/etc/profile -> dotfiles/core
/etc/profile.env -> dotfiles/core
/etc/bash_completion.d -> dotfiles/core
/etc/gnupg -> dotfiles/core
/usr/share/keymaps/i386/qwerty/custom_mapping.map.gz -> dotfiles/core
/home/ari/Documents/custom_mapping.map -> dotfiles/core
/etc/portage -> dotfiles/portage
/usr/src/linux/.config -> dotfiles/linux
/usr/share/themes/mumble-dark.qbtheme -> dotfiles/qbittorrent
/home/ari/Documents/wpa_cli_fix/doc.tex -> dotfiles/fix/wpa_cli_fix.tex
/home/ari/Documents/new_kernel_gentoo/doc.md -> dotfiles/fix/gentoo_new_kernel.md
/home/ari/.racket -> dotfiles/lang
/etc/init.d/kbdrate -> dotfiles/openrc
