# dotfiles/config/systemd/user/kite-updater.service
[Unit]
Description=Kite Updater

[Service]
Type=oneshot
StandardOutput=null
ExecStart=/home/ari/.local/share/kite/current/kite-update "update" "silent"
