# dotfiles/openrc/kbdrate
#!/sbin/openrc-run

description="Applies the rate of the keyboard."

depend()
{
  after keyboard
  keyword -docker -lxc -openvz -prefix -systemd-nspawn -uml -vserver -xenu
}

start()
{
  local figured_kbd_delay="${kbd_delay:-300}"
  local figured_kbd_rate="${kbd_rate:-50}"

  ebegin "Setting the keyboard rate to ${figured_kbd_rate}Hz on ${figured_kbd_delay}ms delay"
  kbdrate --delay="$figured_kbd_delay" --rate="$figured_kbd_rate"
  eend $? "Error key rate to ${figured_kbd_rate}Hz on ${figured_kbd_delay}ms delay" || return $?

  return 0
}
